# dbdk/ai_based_question_recommendation_model/fm-api/config/supervisor/conf.d/flower.conf
; ==========================
;  flower supervisor
; ==========================

[program:flower]
command=/data/fm/venv/fm-api/bin/flower -A worker --address=0.0.0.0 --port=5555 --broker=amqp://guest:guest@localhost:5672//
directory=/data/fm/fm-api
user=fm
numprocs=1
autostart=true
autorestart=true
startsecs=15
stopwaitsecs=15

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true