# dbdk/ai_based_question_recommendation_model/fm-api/config/supervisor/conf.d/fm-api.conf
; ==========================
;  FM-API supervisor
; ==========================

[program:fm-api]
command=/data/fm/venv/fm-api/bin/python /data/fm/fm-api/server.py
directory=/data/fm/fm-api
user=fm
numprocs=1
autostart=true
autorestart=true
startsecs=15
stopwaitsecs=15

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true